# deploy/systemd/dorce-deploy.service
[Unit]
Description=Dorce.ai Git deploy
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/srv/dorce.ai
ExecStart=/bin/bash /srv/dorce.ai/deploy/scripts/deploy.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
